# GPU/TPCCAGPUTracking/tracking-ca/standalone/cuinstall/86-nvidia.rules
SUBSYSTEM=="module", KERNEL=="nvidia", RUN+="/opt/cuda/mkcudadevnodes.sh"
